# Library/LaunchAgents/sh.brew.upgrade.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>sh.brew.upgrade</string>
	<key>Program</key>
	<string>/opt/homebrew/bin/brew</string>
	<key>ProgramArguments</key>
	<string>upgrade</string>
	<key>StandardOutPath</key>
	<string>/tmp/com.rogueserenity.homebrew.upgrade.out</string>
	<key>StandardErrorPath</key>
	<string>/tmp/com.rogueserenity.homebrew.upgrade.err</string>
	<key>StartCalendarInterval</key>
	<array>
		<dict>
			<key>Weekday</key>
			<integer>1</integer>
			<key>Hour</key>
			<integer>7</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
	</array>
</dict>
</plist>
